# docker/services/services/radon-monitor.service
[Unit]
Description=Radon monitor daemon (fills + IB orders)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/app
EnvironmentFile=-/app/.env
EnvironmentFile=-/app/web/.env
ExecStart=/usr/local/bin/python3.13 -m scripts.monitor_daemon.run --once
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
