# docker/services/services/radon-analyst-ratings.service
[Unit]
Description=Radon analyst ratings refresh
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/app
EnvironmentFile=-/app/.env
EnvironmentFile=-/app/web/.env
ExecStart=/usr/local/bin/python3.13 scripts/fetch_analyst_ratings.py --json
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
